# deploy/com.ptzcommander.multiuser.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.ptzcommander.multiuser</string>

    <key>ProgramArguments</key>
    <array>
      <string>__NODE_BIN__</string>
      <string>__ROOT__/dist/index.cjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__ROOT__</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>PORT</key>
      <string>__PORT__</string>
      <key>NODE_ENV</key>
      <string>production</string>
      <key>SESSION_SECRET</key>
      <string>__SESSION_SECRET__</string>
      <key>SESSION_COOKIE_SECURE</key>
      <string>false</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/ptzcommander.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/ptzcommander.stderr.log</string>
  </dict>
</plist>
